# others/systemd/synapse.service
# Systemd service file for Synapse
# Installation:
#   1. Copy this file to /etc/systemd/system/synapse.service
#   2. Update paths and configuration as needed
#   3. Run: systemctl daemon-reload
#   4. Run: systemctl enable synapse
#   5. Run: systemctl start synapse

[Unit]
Description=Synapse - High-performance reverse proxy and firewall
Documentation=https://github.com/gen0sec/synapse
After=network-online.target redis.target
Wants=network-online.target

[Service]
Type=forking
PIDFile=/var/run/synapse.pid

# Paths
ExecStart=/usr/local/bin/synapse --daemon --config /etc/synapse/config.yaml
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID

# Restart policy
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s

# Environment
Environment="RUST_LOG=info"

# Working directory
WorkingDirectory=/var/lib/synapse

# Run as root with full privileges for XDP/BPF
User=root
Group=root

# Security settings - disabled for XDP/BPF to work
# NoNewPrivileges must be false for XDP/BPF operations
NoNewPrivileges=false
PrivateTmp=false

# File system access - minimal restrictions for XDP/BPF
ProtectSystem=false
ProtectHome=false
# XDP/BPF requires access to /sys/fs/bpf and network interfaces
ReadWritePaths=/var/log/synapse /var/run /var/lib/synapse /etc/synapse/certs /sys/fs/bpf
ReadOnlyPaths=/sys/class/net /sys/devices /proc/sys/net

# Capabilities (required for XDP/BPF)
AmbientCapabilities=CAP_NET_ADMIN CAP_SYS_ADMIN CAP_BPF CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_ADMIN CAP_SYS_ADMIN CAP_BPF CAP_NET_BIND_SERVICE

# Network - must have full network access for XDP
PrivateNetwork=false

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=synapse

[Install]
WantedBy=multi-user.target
